# lannes.service
[Unit]
Description=Lannes System - breach, telemetry, health, WebRTC
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lannes-system
Environment=PYTHONUNBUFFERED=1
WorkingDirectory=/home/lannes-system
ExecStart=/usr/bin/python3 /home/lannes-system/lannes-repo/lannes_system_w_webrtc.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
